# deploy/firecracker-netns@.service
[Unit]
Description=Network namespace for Firecracker VM %i

[Service]
Type=oneshot
RemainAfterExit=yes

# Network setup/teardown scripts are generated by VMD and placed on the
# host. They create/destroy the network namespace, veth pair, TAP device,
# and nftables rules for the sandbox.
ExecStart=/usr/local/bin/fc-netns-setup %i
ExecStop=/usr/local/bin/fc-netns-teardown %i
