# deploy/superserve-secretsproxy.service
[Unit]
Description=Superserve secrets proxy
After=network-online.target
Wants=network-online.target
# Start before vmd so the unix socket exists when vmd dials it.
Before=superserve-vmd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/secretsproxy

KillMode=process
Restart=always
RestartSec=2

# Environment
EnvironmentFile=/etc/sandbox/secretsproxy.env

# Tighten the runtime: the daemon only needs network egress and the
# sockets it creates itself. No raw filesystem access required.
ProtectHome=yes
ProtectSystem=strict
NoNewPrivileges=yes
PrivateTmp=yes
RuntimeDirectory=secretsproxy
RuntimeDirectoryMode=0750

LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=secretsproxy

[Install]
WantedBy=multi-user.target
